# testing/zabbix-server-mysql/zabbix-server.initd
#!/sbin/openrc-run

depend() {
	need net
	use mysql
	after firewall
}

start_pre() {
	checkpath -d -m 0775 -o zabbix:zabbix /var/run/zabbix
}
